# scripts/com.muses.webhook.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.muses.webhook</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/path/to/Muses/scripts/webhook-server.js</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/path/to/Muses</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>WEBHOOK_PORT</key>
        <string>9000</string>
        <key>WEBHOOK_SECRET</key>
        <string>your-secret-key-here</string>
        <key>PROJECT_DIR</key>
        <string>/path/to/Muses</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/var/log/muses-webhook.log</string>
    
    <key>StandardErrorPath</key>
    <string>/var/log/muses-webhook-error.log</string>
</dict>
</plist>